# systemd/geospatial-daily-updates.timer
[Unit]
Description=Daily Geospatial Data Update Timer (2:00 AM)
Requires=geospatial-daily-updates.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
